# files/wordpress/generate-secret-token
#!/sbin/openrc-run
CONFIG_FILE=/var/www/localhost/htdocs/wp-config.php
PHRASE="put your unique phrase here"

start () {
	ebegin "Generating secret token for wordpress..."
	while grep -q "$PHRASE" $CONFIG_FILE; do
  		KEY=`cat /dev/urandom | base64 | tr -c -d "[:alnum:]" | head -c 16`
		sed -i '1,/'"$PHRASE"'/s/'"$PHRASE"'/'$KEY'/' $CONFIG_FILE
	done
	eend $?
}
