# files/nextcloud/app-init
#!/sbin/openrc-run

DOCUMENT_ROOT=/var/www/localhost/htdocs
CONFIG_PHP=${DOCUMENT_ROOT}/config/config.php

start() {
    ebegin "Starting ${SVCNAME}"
    if [ ! -f "${DOCUMENT_ROOT}/data/.ocdata" ]; then
        cd "${DOCUMENT_ROOT}"
        occ maintenance:install --database=mysql --database-name=nextcloud --database-user=nextcloud --database-pass='' --admin-user=admin --admin-pass=admin
        occ config:system:delete trusted_domains
        occ config:system:set trusted_domains 1 --value='*'
    fi
    eend 0
}
