# files/walbrix/network-bridge.init
#!/sbin/runscript

depend() {
	need net
}

start() {
	ebegin "Starting" ${SVCNAME}
	/etc/xen/scripts/network-multi-bridge start
	eend $?
}

stop() {
	ebegin "Stopping" ${SVCNAME}
	/etc/xen/scripts/network-multi-bridge stop
	eend $?
}

